# crescent.service
[Unit]
Description=Crescent AGI — Evolving Autonomous Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=faris
WorkingDirectory=/home/faris/crescent-agi
Environment=PATH=/home/faris/crescent-agi/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/faris/crescent-agi/.env
ExecStart=/home/faris/crescent-agi/venv/bin/python main.py
Restart=always
RestartSec=30
StandardOutput=append:/home/faris/crescent-agi/crescent.log
StandardError=append:/home/faris/crescent-agi/crescent.log

# Resource limits
MemoryMax=400M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
